# cli/cloudsack/templates/swift/build_image/supervisord.conf
[supervisord]
nodaemon=true

[program:proxy_server]
command=/usr/bin/python /usr/bin/swift-proxy-server /etc/swift/proxy-server.conf
